# src/etc/crontab
# m h dom mon dow user  command
3 *    * * *   root    /etc/squid3/dnscheck.py

