# src_broker/etc/systemd/system/babeld-monitor.timer
[Unit]
Description=babeld-monitor timer

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target

