# src/etc/systemd/system/babeld.service
[Unit]
Description=babeld for mesh routing
After=network.target auditd.service

[Service]
Type=simple
WorkingDir=/usr/local/bin
ExecStart=/usr/local/bin/babeld -F -S /var/lib/babeld/state -c /etc/babeld.conf
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
